# trading_system/infrastructure/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:apex-api]
command=uvicorn trading_system.main:app --host 0.0.0.0 --port 8000 --workers 2
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/apex-api.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/apex-api-err.log
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app"

[program:apex-dashboard]
command=streamlit run trading_system/dashboard/app.py --server.port 8501 --server.address 0.0.0.0 --server.headless true
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=15
stdout_logfile=/var/log/supervisor/apex-dashboard.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/apex-dashboard-err.log
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app"

[group:apex]
programs=apex-api,apex-dashboard
